# bin/trunkStatus.service
[Unit]
Description=Peer Status Trunk
After=networking.target iptables.service firewalld.service ip6tables.service ipset.service nftables.service asterisk.service

[Service]
Type=simple
ExecStart=/bin/python /opt/trunkstatuscv-4.1.py
RemainAfterExit=yes
Restart=always
RestartPreventExitStatus=0 255


[Install]
WantedBy=multi-user.target

